# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:mihomo]
command=/usr/local/bin/mihomo -d /etc/mihomo
autostart=true
autorestart=true
priority=1
startsecs=3
startretries=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/mihomo.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/mihomo.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:server]
command=node /app/server.js
autostart=true
autorestart=true
priority=2
startsecs=2
startretries=5
environment=SERVER_PORT="9090",MIHOMO_API="http://127.0.0.1:9097",CONFIG_DIR="/etc/mihomo",DATA_DIR="/data"
stdout_logfile=/var/log/supervisor/server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/server.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
